# k8s-ansible_all_in_one_ingress_nginx/roles/nerdctl/files/buildkitd.service
[Unit]
Description=BuildKit Daemon
After=network.target

[Service]
ExecStart=/usr/bin/buildkitd
Restart=always

[Install]
WantedBy=multi-user.target

